# empty_trash/empty_trash.service
[Unit]
Description=empty trash bins

[Service]
Type=oneshot
ExecStart=trash-empty -f