# empty_trash/empty_trash.timer
[Unit]
Description=empty trash bins
[Timer]
OnUnitActiveSec=120m
OnBootSec=10m
[Install]
WantedBy=timers.target