# empty_trash/empty_trash.service
[Unit]
Description=empty trash bin
[Service]
Type=oneshot
ExecStart=trash-empty -f