# backup/restic-backup.timer
[Unit]
Description=Backup with restic daily
[Timer]
OnBootSec=30min
OnUnitActiveSec=2h
[Install]
WantedBy=timers.target