# db/backup/rethinkdb_backup.timer
[Unit]
Description=Backup with rethinkdb daily

[Timer]
OnBootSec=30min
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
